# infra/cleanup/argus-attachment-cleanup.service
# argus — reap expired encrypted attachments (Backblaze B2 + DB). Triggered by the .timer (daily).
# Install: copy this dir to /opt/argus/cleanup, place the units in /etc/systemd/system, then
#   systemctl daemon-reload && systemctl enable --now argus-attachment-cleanup.timer
# On the real deploy, deploy.sh step 5c stages + installs + arms this automatically. Manual dry run:
#   systemctl start argus-attachment-cleanup  (then check `journalctl -u argus-attachment-cleanup`)
# (verify the AWS CLI v2 build + the `docker compose exec` path work under the sandbox as User=argus; MDWE is
# already dropped below.)
[Unit]
Description=argus — reap expired encrypted attachments (Backblaze B2 + DB)
Documentation=https://github.com/Dezoxy/secmes/blob/main/docs/threat-models/encrypted-attachments.md
After=network-online.target
Wants=network-online.target
# The worker runs `docker compose exec` against the postgres container, so it needs the Docker daemon up —
# order after it (and pull it up if idle) so a timer firing after a reboot, before Docker is ready, doesn't
# fail. If Docker is genuinely down the exec fails and OnFailure= alerts.
After=docker.service
Wants=docker.service
# The LoadCredential sources below are delivered by argus-secrets.service (Key Vault → /run/argus/secrets via
# the Managed Identity). Gate on it so a missing secret fails closed instead of starting on a stale file.
Requires=argus-secrets.service
After=argus-secrets.service
# On any failure fire the notifier template (see infra/notify/).
OnFailure=argus-notify-failure@%p.service

[Service]
Type=oneshot
User=argus
Group=argus
TimeoutStartSec=900

# --- Non-secret config (edit for your deployment). The access-key-id is NOT a secret (it rides in every
#     presigned URL); the bucket/endpoint/region are public. ---
# The DB is reached IN-CONTAINER via `docker compose exec` — PG has no published host port (invariant #3),
# so there is no PGHOST/PGPORT/PGPASSFILE. COMPOSE_FILE + COMPOSE_PROJECT_NAME let `docker compose` attach to
# the deployed stack from this timer; PGUSER/PGDATABASE are the role/db the in-container client connects as.
Environment=COMPOSE_FILE=/opt/argus/compose.prod.yaml
Environment=COMPOSE_PROJECT_NAME=argus-prod
Environment=PGUSER=argus_cleanup
Environment=PGDATABASE=argus
Environment=S3_ENDPOINT=https://s3.eu-central-003.backblazeb2.com
Environment=S3_REGION=eu-central-003
# The attachment bucket — templated per deploy from the SAME S3_BUCKET the api uses (cleanup reaps the api's
# own expired attachments), so the worker always targets the configured environment's bucket.
Environment=S3_BUCKET=REPLACE_WITH_ATTACHMENT_BUCKET
# This worker DELETEs from the attachment bucket, so it needs a credential with deleteFiles there. It uses the
# broad `argus-b2-app-key` (key-id = B2_APP_KEY_ID, same as the backup worker) — which spans BOTH buckets with
# delete on both (the BKP-2 over-broad key). It deliberately does NOT use the api's attachment key
# (`argus-s3-secret-access-key`/S3_ACCESS_KEY_ID): that one mints presigned PUT/GET URLs and need not carry
# deleteFiles, so pairing it here would risk a silent 403 on delete. When BKP-2 splits the keys into
# least-privilege bucket-scoped keys, give this worker a credential scoped to the attachment bucket WITH
# delete, and re-point both the key-id and the LoadCredential source below.
Environment=S3_ACCESS_KEY_ID=REPLACE_WITH_B2_KEY_ID

# --- Secret: the B2 application key, delivered as a credential FILE, never in env/Helm at rest (invariant
#     #5). The source path is populated from Azure Key Vault by the Managed Identity at boot (a separate fetch
#     step). systemd copies it into $CREDENTIALS_DIRECTORY (mode 0400, tmpfs), referenced as %d below. (No DB
#     password: the DB connection runs in-container over the local-trust socket — see cleanup-attachments.sh.) ---
LoadCredential=b2-secret:/run/argus/secrets/b2-app-key
Environment=S3_SECRET_ACCESS_KEY_FILE=%d/b2-secret

# AWS CLI v2 wants a writable HOME/cache, but ProtectHome + ProtectSystem make the FS read-only.
# Give it a private tmpfs runtime dir as HOME so it never touches the real home.
RuntimeDirectory=argus-cleanup
# 0700 (not the 0755 default): the worker materialises the AWS credentials file here, on this
# service-private tmpfs (no host-disk backing) — so other users can't even list it.
RuntimeDirectoryMode=0700
Environment=HOME=%t/argus-cleanup
# `docker compose` needs the staged compose tree + a HOME to resolve from; run from the app dir.
WorkingDirectory=/opt/argus

ExecStart=/opt/argus/cleanup/cleanup-attachments.sh

# --- Hardening: least-privilege oneshot that reaches the docker socket (to exec psql in-container) + B2
#     egress. `argus` is already in the docker group (cloud-init), so this grants no NEW privilege. ---
NoNewPrivileges=true
# ProtectSystem=full (not strict): the docker client must reach /run/docker.sock; strict would make /run
# read-only and block the connect. full still keeps /usr,/boot,/etc read-only; CapabilityBoundingSet stays empty.
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectHostname=true
ProtectClock=true
ProtectProc=invisible
ProcSubset=pid
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
CapabilityBoundingSet=
LockPersonality=true
# MemoryDenyWriteExecute is deliberately NOT set: this worker runs AWS CLI v2 (a PyInstaller bundle that maps
# memory W+X), which MDWE would kill — see argus-db-backup.service for the full rationale. ~18 other
# hardening directives still apply.
SystemCallArchitectures=native
UMask=0077
