# infra/stack/secrets/argus-secrets.service
# argus — fetch runtime secrets from Key Vault (Managed Identity) → credential files on tmpfs.
# Install: copy this dir to /opt/argus/secrets, place this unit in /etc/systemd/system, set ARGUS_KEY_VAULT
# below (the Slice-4 deploy templates it from the `key_vault_name` Terraform output), then
#   systemctl daemon-reload && systemctl enable --now argus-secrets.service
# Verify (the journal logs secret NAMES + status only, never values):  systemctl start argus-secrets && \
#   systemctl status argus-secrets
#
# Consumers MUST order after + require this unit so they never start on a missing/stale secret, e.g.:
#   [Unit]
#   Requires=argus-secrets.service
#   After=argus-secrets.service
# (the Compose stack launcher + argus-db-backup + argus-attachment-cleanup all read /run/argus/secrets/*).
[Unit]
Description=argus — fetch Key Vault secrets via Managed Identity to /run/argus/secrets
Documentation=https://github.com/Dezoxy/secmes/blob/main/infra/stack/secrets/README.md
After=network-online.target
Wants=network-online.target
# When ARGUS_TOKEN_SOURCE=arc (the AWS/EC2 experiment), the token comes from the Azure Arc agent's local HIMDS
# endpoint, so order AFTER its daemon (himdsd.service — the Connected Machine Agent's HIMDS service). WEAK
# dependency (Wants/After, not Requires): the live Azure VM has no himdsd.service, where Requires would fail the
# unit — Wants is a no-op when the target is absent, and the fetch script's curl --retry absorbs a slow agent.
# On EC2 the Arc agent is present and this orders correctly.
After=himdsd.service
Wants=himdsd.service
# Secrets must exist before anything that consumes them comes up.
Before=docker.service

[Service]
Type=oneshot
# Stay "active" after exit so dependents (Requires=) see success and can start.
RemainAfterExit=yes
TimeoutStartSec=120

# Non-secret deploy config. ARGUS_KEY_VAULT is the Key Vault NAME (the `key_vault_name` TF output); the token
# itself comes from the VM Managed Identity via IMDS — there is no static credential here. Replace on deploy.
Environment=ARGUS_KEY_VAULT=REPLACE_WITH_KEY_VAULT_NAME
Environment=ARGUS_SECRETS_DIR=/run/argus/secrets
ExecStart=/opt/argus/secrets/fetch-keyvault-secrets.sh

# --- Hardening: a short-lived root oneshot that needs only egress (IMDS + Key Vault) and a write to the
#     secrets tmpfs. Runs as root because it writes 0400 root:root files the Docker daemon + systemd
#     LoadCredential (both root) read. curl + jq only — no Python — so MemoryDenyWriteExecute is safe. ---
NoNewPrivileges=true
ProtectSystem=strict
# Let systemd create /run/argus/secrets (root:root, 0700) BEFORE ExecStart — so it survives a reboot that
# wipes /run, and it's writable without a ReadWritePaths entry (a ReadWritePaths on a missing path fails the
# unit during namespace setup, before the script could mkdir it). Preserve=yes keeps the dir + its contents
# across restarts (the stack + backup/cleanup units read it and this unit is RemainAfterExit).
RuntimeDirectory=argus/secrets
RuntimeDirectoryMode=0700
RuntimeDirectoryPreserve=yes
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectKernelLogs=true
ProtectHostname=true
ProtectClock=true
ProtectProc=invisible
ProcSubset=pid
# AF_UNIX is needed for NSS/DNS resolution of the vault hostname (systemd-resolved/nscd talk over a UNIX
# socket); matches the sibling backup/cleanup units. IMDS is an IP literal, but Key Vault is a DNS name.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
CapabilityBoundingSet=
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
UMask=0077

[Install]
WantedBy=multi-user.target
