# infra/backup/argus-db-backup.timer
# argus — daily timer for the encrypted DB backup worker. Runs in the early-morning quiet window (UTC).
# Storage stays bounded via a server-side B2 lifecycle rule (Object Lock / BKP-2), not a script prune — the
# worker's key has no delete capability.
[Unit]
Description=argus — nightly encrypted DB backup

[Timer]
# 02:30 UTC daily. Persistent=true runs a missed backup after downtime (a skipped night is a real gap).
OnCalendar=*-*-* 02:30:00
Persistent=true
RandomizedDelaySec=10min

[Install]
WantedBy=timers.target
