# infra/backup/argus-db-backup.service
# argus — nightly logical DB backup (pg_dump → age → Backblaze B2). Triggered by the .timer (daily).
# Install: copy this dir to /opt/argus/backup, place the units in /etc/systemd/system, then
#   systemctl daemon-reload && systemctl enable --now argus-db-backup.timer
# On the real deploy, deploy.sh step 5c stages + installs + arms this automatically. Manual dry run:
#   systemctl start argus-db-backup  (then check `journalctl -u argus-db-backup`)
# (verify the AWS CLI v2 build + the `docker compose exec` path work under the sandbox as User=argus; MDWE is
# already dropped below — see README's "Deploy verification").
# Prerequisites: the argus_backup role provisioned LOGIN (no password — in-container local trust; see README),
# and an age keypair whose PUBLIC key is set in AGE_RECIPIENT below (the PRIVATE key stays in Key Vault).
[Unit]
Description=argus — nightly encrypted DB backup (pg_dump → age → B2)
Documentation=https://github.com/Dezoxy/secmes/blob/main/infra/backup/README.md
After=network-online.target
Wants=network-online.target
# The LoadCredential sources below are delivered by argus-secrets.service (Key Vault → /run/argus/secrets via
# the Managed Identity). Gate on it so a missing secret fails closed instead of starting on a stale file.
Requires=argus-secrets.service
After=argus-secrets.service
# On any failure (non-zero exit, timeout, signal) fire the notifier template, which posts a Sentry event
# to GlitchTip. %p = unit prefix without .service, so the instance is argus-notify-failure@argus-db-backup.
OnFailure=argus-notify-failure@%p.service

[Service]
Type=oneshot
User=argus
Group=argus
# A full dump + upload can take a while as the DB grows; cap generously.
TimeoutStartSec=3600

# --- Non-secret config (edit for your deployment). The access-key-id is NOT a secret (it rides in every
#     presigned URL); the bucket/endpoint/region and the age PUBLIC recipient key are public. ---
# The DB is reached IN-CONTAINER via `docker compose exec` — PG has no published host port (invariant #3),
# so there is no PGHOST/PGPORT/PGPASSFILE. COMPOSE_FILE + COMPOSE_PROJECT_NAME let `docker compose` attach to
# the deployed stack from this timer; PGUSER/PGDATABASE are the role/db the in-container client connects as.
Environment=COMPOSE_FILE=/opt/argus/compose.prod.yaml
Environment=COMPOSE_PROJECT_NAME=argus-prod
Environment=PGUSER=argus_backup
Environment=PGDATABASE=argus
Environment=S3_ENDPOINT=https://s3.eu-central-003.backblazeb2.com
Environment=S3_REGION=eu-central-003
# The PRIVATE db-backups bucket — separate from the attachment bucket, never public. Templated per deploy
# from BACKUP_S3_BUCKET so the AWS experiment can't write its dumps into the production backup bucket.
Environment=S3_BUCKET=REPLACE_WITH_BACKUP_BUCKET
Environment=S3_ACCESS_KEY_ID=REPLACE_WITH_B2_KEY_ID
# age PUBLIC recipient key (age1...). The matching PRIVATE key is NOT here — it lives in Key Vault and is
# fetched only at restore time, so this host cannot decrypt the backups it writes.
Environment=AGE_RECIPIENT=REPLACE_WITH_AGE_PUBLIC_KEY
Environment=RETENTION_DAYS=30

# --- Secret: the B2 application key, delivered as a credential FILE, never in env/Helm at rest (invariant
#     #5). The source path is populated from Azure Key Vault by the Managed Identity at boot (a separate fetch
#     step). systemd copies it into $CREDENTIALS_DIRECTORY (mode 0400, tmpfs), referenced as %d below. (No DB
#     password: the DB connection runs in-container over the local-trust socket — see backup-db.sh.) ---
LoadCredential=b2-secret:/run/argus/secrets/b2-app-key
Environment=S3_SECRET_ACCESS_KEY_FILE=%d/b2-secret

# AWS CLI v2 wants a writable HOME/cache, but ProtectHome + ProtectSystem make the FS read-only.
# Give it a private tmpfs runtime dir as HOME so it never touches the real home.
RuntimeDirectory=argus-db-backup
# 0700 (not the 0755 default): the worker materialises the AWS credentials file here, on this
# service-private tmpfs (no host-disk backing) — so other users can't even list it.
RuntimeDirectoryMode=0700
Environment=HOME=%t/argus-db-backup
# `docker compose` needs the staged compose tree + a HOME to resolve from; run from the app dir.
WorkingDirectory=/opt/argus

ExecStart=/opt/argus/backup/backup-db.sh

# --- Hardening: least-privilege oneshot that reaches the docker socket (to exec pg_dump in-container) + B2
#     egress. `argus` is already in the docker group (cloud-init), so this grants no NEW privilege. ---
NoNewPrivileges=true
# ProtectSystem=full (not strict): the docker client must reach /run/docker.sock; strict would make /run
# read-only and block the connect. full still keeps /usr,/boot,/etc read-only. The privileged work is the
# daemon's, not this unit's (its CapabilityBoundingSet stays empty).
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectHostname=true
ProtectClock=true
ProtectProc=invisible
ProcSubset=pid
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
CapabilityBoundingSet=
LockPersonality=true
# MemoryDenyWriteExecute is deliberately NOT set: AWS CLI v2 ships as a PyInstaller bundle whose bootloader
# maps embedded modules W+X, which MDWE kills — the worker would segfault and the backup silently fail
# (BKP-1). The README documents this as the one knob that bites the CLI; the other ~18 directives still apply.
SystemCallArchitectures=native
UMask=0077
